# OpenWrt/etc/init.d/gas_counter
#!/bin/sh /etc/rc.common
# https://github.com/tlinnet/gas-counter-magnetic/tree/main

START=95
STOP=15

start() {
echo start

while true  # Keep an infinite loop to reconnect when connection lost/broker unavailable
do
    mosquitto_sub -h slateplus.lan  -u "gasread" -P "Hello" -v -t "sensors/gas/#" --qos 1 --id "gasread" --disable-clean-session | while read -r payload
    do
        # Here is the callback to execute whenever you receive a message:
        echo "Rx MQTT: ${payload}" | tee -a /root/gas_data.csv
    done
    sleep 10  # Wait 10 seconds until reconnection
done # &  # Discomment the & to run in background (but you should rather run THIS script in background)

}

stop() {
echo stop
}
