# OpenWrt/etc/init.d/gas_counter
#!/bin/sh /etc/rc.common
# https://github.com/tlinnet/gas-counter-magnetic/tree/main

START=95
STOP=15

start() {
echo start
mosquitto_sub -h slateplus.lan  -u "gasread" -P "Hello" -v -t "sensors/gas/#" --qos 1 --id "gasread" --disable-clean-session |
while read payload; do
  echo Received $payload | tee -a /root/gas_data.csv
done
}

stop() {
echo stop
}
