# OpenWrt/etc/init.d/gas_counter
#!/bin/sh /etc/rc.common
# https://github.com/tlinnet/gas-counter-magnetic/tree/main

START=95
STOP=15

start() {
        echo start
        # commands to launch application
        mosquitto_sub -h slateplus.lan  -u "gasread" -P "Hello" -v -t "sensors/gas/#" --qos 1 --id "gasread" --disable-clean-session  | tee /root/gas_data.csv
}

stop() {
        echo stop
        # commands to kill application
}
